# contrib/findmissing/config/chromeos-patches-app.conf
; Copyright 2020 The Chromium OS Authors. All rights reserved.
; Use of this source code is governed by a BSD-style license that can be
; found in the LICENSE file.

[program:chromeos_patches]
directory=/home/chromeos_patches
command=/home/chromeos_patches/env/bin/honcho start -f ./config/procfile chromeos_patches_server
autostart=true
autorestart=true
user=chromeos_patches

# Environment variables ensure that the application runs inside of the configured virtualenv.
environment=VIRTUAL_ENV="/home/chromeos_patches/env",PATH="/home/chromeos_patches/env/bin",HOME="/home/chromeos_patches",USER="chromeos_patches"
stdout_logfile=syslog
stderr_logfile=syslog
